# deploy/systemd/wp-handover.timer
[Unit]
Description=wp-tt authoritative soon-idle labels every 60s (incremental JOB_ORDER_HISTORY completions)

[Timer]
OnBootSec=2min
OnUnitActiveSec=60s
AccuracySec=10s

[Install]
WantedBy=timers.target
