# deploy/systemd/wp-nightly.service
[Unit]
Description=wp-tt-dashboard nightly authoritative KPI extract (yesterday, all KPIs)
After=network-online.target

[Service]
Type=oneshot
EnvironmentFile=%h/projects/wp-tt-dashboard/.env
# runs all 8 extracts for yesterday + transform (final) + baseline
ExecStart=%h/projects/wp-tt-dashboard/target/release/extractor run --kpi all --target oracle-prod
Nice=10
