# deploy/systemd/wp-shift-t2.service
[Unit]
Description=wp-tt-dashboard LIVE shift tick T2 (heavier JOB_ORDER_HISTORY: empty, cycle, crane-q)
After=network-online.target

[Service]
Type=oneshot
EnvironmentFile=%h/projects/wp-tt-dashboard/.env
ExecStart=%h/projects/wp-tt-dashboard/target/release/extractor tick --shift --tier t2 --target oracle-prod
Nice=18
