# deploy/systemd/wp-shift-t1.service
[Unit]
Description=wp-tt-dashboard LIVE shift tick T1 (cheap: MPH, QC-wait, util + vessels)
After=network-online.target

[Service]
Type=oneshot
EnvironmentFile=%h/projects/wp-tt-dashboard/.env
ExecStart=%h/projects/wp-tt-dashboard/target/release/extractor tick --shift --tier t1 --target oracle-prod
Nice=15
