# deploy/systemd/wp-tick-t2.timer
[Unit]
Description=wp-tt tick T2 every 20 min (load-conscious; heavier source)

[Timer]
OnBootSec=5min
OnUnitActiveSec=20min

[Install]
WantedBy=timers.target
