# deploy/systemd/wp-tick-t2.service
[Unit]
Description=wp-tt-dashboard intra-day tick T2 (heavier JOB_ORDER_HISTORY KPIs, today provisional)

[Service]
Type=oneshot
EnvironmentFile=%h/projects/wp-tt-dashboard/.env
ExecStart=%h/projects/wp-tt-dashboard/target/release/extractor tick --tier t2 --target oracle-prod
Nice=15
