# deploy/systemd/wp-workpool.timer
[Unit]
Description=wp-tt LIVE work-pool snapshot every 90s (one bounded JOB_QUEUE_SCHEDULE + JOB_ORDER_LIST pull)

[Timer]
OnBootSec=1min
OnUnitActiveSec=90s
AccuracySec=10s

[Install]
WantedBy=timers.target
