# deploy/systemd/wp-workpool.service
[Unit]
Description=wp-tt-dashboard LIVE work-pool snapshot (JOB_QUEUE_SCHEDULE + JOB_ORDER_LIST -> Postgres)
After=network-online.target

[Service]
Type=oneshot
EnvironmentFile=%h/projects/wp-tt-dashboard/.env
ExecStart=%h/projects/wp-tt-dashboard/target/release/extractor workpool --target oracle-prod
Nice=15
